# deploy/ib-bots.service
[Unit]
Description=IB Trader Bot Runner
After=ib-engine.service
Requires=ib-engine.service

[Service]
Type=simple
User=zuura
WorkingDirectory=/home/zuura/projects/ib-trader
ExecStart=/home/zuura/projects/ib-trader/.venv/bin/ib-bots
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
